# deploy/systemd/tranzit-api-docker.service
[Unit]
Description=Tranzit API (Docker)
Documentation=https://github.com/kkzakaria/tranzit
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
Restart=on-failure
RestartSec=5
SuccessExitStatus=143 137
TimeoutStartSec=30
TimeoutStopSec=15
# KillMode=none is intentional: docker run runs in the foreground and systemd
# must not kill the docker CLI process directly — stop is handled entirely by
# ExecStop/ExecStopPost. Generates a deprecation warning on systemd 246+; safe
# to ignore for this Docker-managed lifecycle pattern.
KillMode=none

ExecStartPre=-/usr/bin/docker stop tranzit-api
ExecStartPre=-/usr/bin/docker rm tranzit-api
ExecStart=/usr/bin/docker run \
  --name tranzit-api \
  --env-file /etc/tranzit/api.env \
  -p 34001:34001 \
  --log-driver=journald \
  --log-opt tag=tranzit-api \
  tranzit-api:latest
ExecStop=/usr/bin/docker stop -t 10 tranzit-api
ExecStopPost=-/usr/bin/docker rm -f tranzit-api

# Container logs go directly to journald via --log-driver=journald.
# Suppress stdout/stderr here to avoid double-logging from the docker CLI wrapper.
# To follow logs: journalctl -f CONTAINER_NAME=tranzit-api
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
